# crontab
# Executa o comando Django para realizar o backup diariamente às 3h00
# 0 3 * * * root python /code/manage.py dbbackup
20 13 * * * root python /code/manage.py dbbackup
